# git/.config/profile.d/20-git.sh#template
#!/bin/sh

export GITHUB_TOKEN="{{ required .Vars.github.cli_token }}"
