# git/.config/profile.d/20-git.sh#template
#!/bin/sh

export GITHUB_TOKEN="{{ bitwarden "github.com" "field" "gh-cli-token" }}"
