# env/.config/profile.d/02-env-secret.sh#encrypted
-----BEGIN AGE ENCRYPTED FILE-----
YWdlLWVuY3J5cHRpb24ub3JnL3YxCi0+IFgyNTUxOSBMcVpzUy8rWTIvRE1peFR4
MTRzcmFKSXNxdGx5UTI0cjJVRmc0eUowaUZBCkZaeUswR1VrZUFjUFQreUdWak9w
Qld1UEZ3NmQxcE1uYTVEc2NRUE40alEKLS0tIEk4c1gwMEJvS0RDWlh4VzBFSUp4
R0hQeFEvb3VqenJXeUpRcXBYcWJsTUUKp33MgpaGwHTsTSawSebX2dMTag0+8z49
QH8FOfqbqp1qLiUTMd9/zHlFTuB+LqHypdn7qWCEu6qde/h/vnSFpziUOd+fFYKq
ycc=
-----END AGE ENCRYPTED FILE-----
