# env/.config/profile.d/10-base16.sh
#!/bin/sh

BASE16_COLOR_00_HEX=${BASE16_COLOR_00_HEX:-#262626}
BASE16_COLOR_01_HEX=${BASE16_COLOR_01_HEX:-#d75f5f}
BASE16_COLOR_02_HEX=${BASE16_COLOR_02_HEX:-#afffaf}
BASE16_COLOR_03_HEX=${BASE16_COLOR_03_HEX:-#d7d7af}
BASE16_COLOR_04_HEX=${BASE16_COLOR_04_HEX:-#87afd7}
BASE16_COLOR_05_HEX=${BASE16_COLOR_05_HEX:-#d7afd7}
BASE16_COLOR_06_HEX=${BASE16_COLOR_06_HEX:-#87d7af}
BASE16_COLOR_07_HEX=${BASE16_COLOR_07_HEX:-#d7d7d7}
BASE16_COLOR_08_HEX=${BASE16_COLOR_08_HEX:-#6c6c6c}
BASE16_COLOR_09_HEX=${BASE16_COLOR_09_HEX:-#ff5f5f}
BASE16_COLOR_0A_HEX=${BASE16_COLOR_0A_HEX:-#afd787}
BASE16_COLOR_0B_HEX=${BASE16_COLOR_0B_HEX:-#ffffd7}
BASE16_COLOR_0C_HEX=${BASE16_COLOR_0C_HEX:-#87d7ff}
BASE16_COLOR_0D_HEX=${BASE16_COLOR_0D_HEX:-#d787ff}
BASE16_COLOR_0E_HEX=${BASE16_COLOR_0E_HEX:-#87ffff}
BASE16_COLOR_0F_HEX=${BASE16_COLOR_0F_HEX:-#d7d7d7}

export BASE16_COLOR_00_HEX BASE16_COLOR_01_HEX BASE16_COLOR_02_HEX BASE16_COLOR_03_HEX \
	BASE16_COLOR_04_HEX BASE16_COLOR_05_HEX BASE16_COLOR_06_HEX BASE16_COLOR_07_HEX \
	BASE16_COLOR_08_HEX BASE16_COLOR_09_HEX BASE16_COLOR_0A_HEX BASE16_COLOR_0B_HEX \
	BASE16_COLOR_0C_HEX BASE16_COLOR_0D_HEX BASE16_COLOR_0E_HEX BASE16_COLOR_0F_HEX
