# env/.config/profile.d/11-aliases-macos.sh#profile:macos
#!/bin/sh

# clipboard
alias cc="fc -ln -1 | pbcopy"
alias xc="pbcopy"
alias xp="pwd | tr -d '\n' | pbcopy"

# misc
alias bu="brew update && brew upgrade --quiet --force"
