# env/.config/profile.d/20-go.sh
#!/bin/sh

export GOPATH="$XDG_DATA_HOME/go"
export GOBIN="$HOME/.local/bin"

prepend_path "$GOPATH/bin"
