# env/.config/profile.d/20-nodejs.sh
#!/bin/sh

export NODE_NO_WARNINGS=1
export NODE_REPL_HISTORY="$XDG_DATA_HOME"/node_repl_history
