# env/.config/profile.d/30-less.sh
#!/bin/sh

export LESS="-RF --mouse"
export LESSHISTFILE="$XDG_STATE_HOME/lesshst"
