# env/.config/profile.d/30-man.sh
#!/bin/sh

export MANPAGER="less -R --use-color -Dd+g -Du+b"
export MANROFFOPT="-P -c"
