# hyprland/.config/profile.d/30-hyprland.sh
#!/bin/sh

alias cc="fc -ln -1 | wl-copy -n"
alias xc="wl-copy -n"
alias xp="pwd | wl-copy -n"
