# bitwarden/.config/profile.d/30-bitwarden.sh
#!/bin/sh

export BITWARDENCLI_APPDATA_DIR="${BITWARDENCLI_APPDATA_DIR:-"$HOME/.local/share/bitwardencli"}"
alias bw=bww
