# eza/.config/profile.d/70-eza.sh
#!/bin/sh

alias l="eza --group-directories-first --group -1"
alias ll="eza --long --all --group-directories-first --group --icons --git"
alias lt="eza --long --all --sort=modified --group-directories-first --group --icons --git"
alias ltree="eza --long --all --sort=modified --group-directories-first --group --icons --git --tree --level=2"
