# atuin/.config/profile.d/50-atuin.sh
#!/bin/sh

export ATUIN_SYNC_ADDRESS=https://atuin.sbbh.cloud
