# atuin/.config/profile.d/70-atuin.sh
#!/bin/sh

ATUIN_SYNC_ADDRESS=https://atuin.sbbh.cloud
