# root_xps/etc/systemd/system/mnt-backup.mount
[Unit]
Description=mount backup NFS share

[Mount]
What=nas.int.sbbh.cloud:/volume1/backups
Where=/mnt/backup
Type=nfs
Options=defaults
TimeoutSec=5

[Install]
WantedBy=multi-user.target
