# root_xps/etc/systemd/system/mnt-backup.automount
[Unit]
Description=automount backup NFS share

[Automount]
Where=/mnt/backup
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
