# root_xps/etc/systemd/system/mnt-workspace.automount
[Unit]
Description=automount workspace NFS share

[Automount]
Where=/mnt/workspace
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
