# tmux/.config/profile.d/20-tmux.sh
#!/bin/sh

alias mux=__tmuxinator_sessions
