# nvim/.config/profile.d/30-nvim.sh
#!/bin/sh

export NVIM_APPNAME=nvim
export NVIM_HOME="$HOME/.config/$NVIM_APPNAME"
export MANPAGER='nvim +Man!'

alias vi=nvim
