# work/.config/profile.d/05-work.sh#template#encrypted
-----BEGIN AGE ENCRYPTED FILE-----
YWdlLWVuY3J5cHRpb24ub3JnL3YxCi0+IFgyNTUxOSBaY0Z3eHVYSHpmQVgrK1Fo
OUZYNjVuckZQNUFTREgwM3dndm5QdFJtL1NZCklNdHVQTVRSbVVaaERNR2R5dDk4
MWJzeGFVOFo3NnB5bWRwcVgya1lVYmcKLS0tIEhlTXBJaHRpb2JPdDVGcTFER1ZL
QXlUY3ZFTGQ4M05hV0RLb2xOZFJHUXcKgBw+BjASlz/x/Ojx6XVrrtB3xSGrhE1V
YhkGYg7jct+pbwcJhMsFmoR8fKApzAZYPggE33/g/tsqtHaSh15KBfgONL8iKHUL
ljpG5mx/m2f0BO4kmRW11WEW1mshIlESBN/I7l2C7t0dlDFD5PhKxGS2eob5c2W9
/FK9LCOwopZfaDUW4CD2s32547RYHNzeswjh/ti2KBEuPbFpaWatnUt8xfzRBkGr
Hs2URMFfiXUHSW23MNNcI1y7JkLFAUhyb+dP5Sp8z7HAbw9D6+M0TsbAf7iBImZj
cb+zbbs+rrUGOYDnCMOWDdAgtM7ehcbxEZuO+cuwRlDjiO9nOuporJOex2nO2gVy
W8XpA6dfwHERZQXnjaNUSWyuAVYvPAEYzQQQ4KRthXojBVItO/9vFOt2/nHL9N/F
DN2vf/5z9sL+wuP2BZ4bficcCj/bxCi5zfhk2U9bJES4PYk/eWDuots/wPLUY9i/
r5UscluwntdbHIqe464wac05hSRpkRHyl9RfBmgIHTZtiedcAwS2aVa/I8W9i95e
Hzl+pZGHQ0CvSaEofuohPMfpZTX7WZanDMBo0cDFaX0a6CJSo6j4acfcOeAZhKGL
sophYRcaKAvEInbYa6ICk9xNjxVHXGy9PAg1zIbgZ7C11B1cp3OwbTeSo4s+WGuq
nrRshh3PuTolJJFXOx7xR3J1NRtINQ0MQlTGqjIG2+yjmHuGPUiEy66L1BQGWIy2
xughaX9G4A0MpJLeDKJr1WewMYt3uAHOgKDTZAzDuquo4HJBDbg+WQNS3zZePUc6
6S3cS28FcTY4PqXsFFQKL4oJQZKGaLUoOWBe573TgmlwxmZhB4zan4HVGfys7dre
WBY4GpywLAKovu1BvT2TezqrJK1bdsC0UMSNTOrzXw2qjCgMWjihSzt5hp+dSZsi
M7eMqpdgBshtBD1+sn9O6ElcohzwrEAMuF0k/6Ixep/jTTidEcmksnrPKj25y/LV
cVbn1yJP/jtPgS0voiyXuo5wZtP2Xxx3h41daJ5va0/qCQw+tG1SY6pGLLaLJdTz
QQu150cfl5hzImprguR7D1PjuMZehZLQYe42GqHSPQsWZNqxgOxKYuijlxQLxYPD
FCbn2RDnR7OZd9Hacy0laPspTAd/1eTsi9olfx1qQ310QLn5HuwG7++TQ7+AIrv3
37W99dzrvTIGA40/mm9wyJNatoK8mnphzcfqGnK7YG882Pf0jelQQ5zOCYbaUNsK
ADB6TEysinyMtPlb99EiLEdN6cZTAZvNV+3XNEjbEw1Hc7oiq+rI6vrMhQbQNPEu
actmkCGS4NlgHpC80v3w11Kxi9jC4CFvoMcRZmFbyrZV5rm8Np6ouuKuI3zF6nAj
+Pq8W/SOE2jtrXVnS2mrmNy55RgC5ZkfoTiReoe3lfIGnityDnhQxcRbVNwMSzHd
yV2HHIrMRq3L0Ztc+HXW3Krs4xzs2ZC7h9IOmpGAk6PxzKDE6ZHoekjaLL2tXmXj
sXRLSE4SNgGDy7cNh1tdWBUA6BYJfuYDDA8KmIOaylb6Cha0krLe7TrvWarQ6h0W
kB3YtJTpPbxVj3YGoqnoP6tHeHvrydWowE2JNhuqaYyjVD4elDVyceZi/tcgicEp
iUtrKHJ0A4bF4cai1OJF5nwzf3FXDanPO+1LCoR5DEqAOFjq4cH9CQ==
-----END AGE ENCRYPTED FILE-----
