# work/.config/profile.d/05-work.sh#template#encrypted
-----BEGIN AGE ENCRYPTED FILE-----
YWdlLWVuY3J5cHRpb24ub3JnL3YxCi0+IFgyNTUxOSBsRWtmbEpCUDFRM1hXU3Bx
YzJIdzlsb0FubDNRNXB4SzRURU93bUQ5K1JFCmdWSUNnUE91U0dyek0rZm82VTNj
d3c1RjF6Z3dFODJ3MVRZLy95bkVFSlUKLS0tIHFqN3BaTUswR2liWno4YkFTS09R
OEpHZFM3VVFGWUtmSGtFSGdxNUxDK00KRSRRmF1u5K09B2ZbtrrEGPuZdWiSp0zy
mS9aanOkUlrHHc9toSmYE9KsOAPhT5speU9U/l78W0tnW99z9MJCP8npRPfmPwcA
SdCht76NJ4ob8LLU1gFQFejYwwM72fwOOYJedBVKJFGyfZnRNRv/JhnWaTOcmruc
WrwMyOs9QnyHcUjM5uE1QUvhKEkDtkfdfJ/A6ECUSvJqP38p1bjp9Z+qp9u4Zinj
xtYoI9Ax8hirRe0xrAMhQrco2XQnRLCrsWy+dxjhU2Fng4XtdCX7tEggXd5ISmje
frT5zczCLcrgOKTY+xSMgSk9ZEwGzGRLLDQBxISKgGD0huXLBv0hfagcitIFtnlY
km3u2yxefDjVwOU4nqcNMWzI0Fo5dkeuTAY7IUHK3rAKYRP7pcxU5pp+gGcnM4c0
iUaoVqzRNDS2Mk9uJmSWrmmWci3KmF0WwYlnEfK5HzTw8HDznEpBj53elDhWrpYS
H+xlvMyYSXbmZqQVzXG4HNJ+4Dr1dIUdAdcKWYGk5wVjc2P4xgMaz7fTtofMNiks
VdWtMXPoMUyxv9SFwI5xMl7mZEJLDhYg3KxkNnwmLQ3Gbsum4VVVgIl2Q6j2N6oH
KKx3H+VVsHqRQjwr6MKK2pRuk26BeCB3x3KXUGnz4pqX9gKOZ2v5nlhzNNfadGy5
IEJMuNooIfKK+Un/OISvb65yAwEFxf6RMq6bji6QvBOkVgi2b5FLHwDcuk483gbC
ZwLOMI9hpVkUfjcSZ68N+WTUXFH0hcx77F2flFgqJhydIwr28X3XFZPisqbfTJPu
Gi7xDq2YGega8UM9PNmX376jU8l7rkS0epzxrgZk1bCFRy8fKwVogdzL1xQw6VpK
JRr2TEygJcMizNn4V4sc4/Q2a1m1r1LeHZCJ8Wce9NRPddew3nfjJ2eekLzi2dIN
R1wI/jxpnyUj8Es7k0cJqWBKI58QNZEIkXnRuJzA1oN0h39CFIQ1SYLAO6Z5kFVD
Gfwg1Ayswj43ApMNyM2+33jeBxcT+QYyGN7nd4hh/npwSUme8LB6VAbdteGiaE/9
nE5KsK2LHyineDCqSmUxZPd2yTBH/TJBDwq7ZJQT/OUeNSapIhtASmlYJYcpoLUl
C+Z3LzCEbpU7olHP0xw4nA3jNNb0qKOELVMrF3XkLrGnd70T+rXuyECrrscy9Y87
acgFazKJYZ+nIjV3EuFyMjkjz/ZLAuEiO8UN1+TAAsVeTc+rD7Baki7Om3O6ZfKK
R+jPxNltOG9KZEeASbyeVlDZMcVdgsUMQaA5/fdtdZV1x8ID/3A1NxZFV6o7VTMV
4lCUNXrxuH0tFUWTS9HMNJ88hK1INUWVZl6JMEZjpFD5OY67f1VA60UBWeKvODZR
i+bB7/PvgJJXS2kysj41nPz2iftsxZZ2XPEKW+yGPpRNi3fXu+KpwBuieIKNZ8AX
/yC8MTja8Dmy5eR0+L/m5UQjE+88OBbO6QxPJGsK4N0RybvCfGhWixMDXjXMgwon
i81+BpPwZuCKsLN8PYlugzT+tGZc1xVJsTUOm/0hMF/Gk2nXLbnKkHKAqczVOT2M
FghFBOTGAZwm4e88SOnhovtaa2QtZZCN98OXlOh4fk8gvEE0jmQs
-----END AGE ENCRYPTED FILE-----
