# starship/.config/profile.d/50-starship.sh
#!/bin/sh

export STARSHIP_LOG="error"
export STARSHIP_CONFIG="${XDG_CONFIG_HOME}/starship/starship.toml"
