# statemate/.config/profile.d/50-statemate.sh
#!/bin/sh

alias cd.="cd \${XDG_DOTFILES_HOME}"
alias ds='mate status'
alias da='mate apply'
alias dad='mate apply --dry-run -v'
alias de='mate edit'
alias ddiff='mate diff'
alias ddoc='mate doctor'
