# root_xps/etc/udev/rules.d/99-x11-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/home/sysadm/.local/bin/bspwm_setup_monitors"
