# hyprland/.config/systemd/user/waybar.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/waybar -c %h/.cache/waybar/config
